# nghttp2/cmake-build-relwithdebinfo/contrib/nghttpx-upstart.conf
# vim: ft=upstart:

description "HTTP/2 reverse proxy"

start on runlevel [2]
stop on runlevel [016]

exec C:/Program Files (x86)/nghttp2/bin/nghttpx
